# infra/ansible/roles/healthcoach/templates/health-sync.timer.j2
[Unit]
Description=Run Health Sync every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
