# infra/ansible/roles/healthcoach/templates/healthcoach-decision.timer.j2
[Unit]
Description=Run Health Coach Decision Engine Daily at 0510 local time
Requires=healthcoach-decision.service

[Timer]
OnCalendar=*-*-* 05:10:00
Persistent=true

[Install]
WantedBy=timers.target
