# infra/ansible/roles/healthcoach/templates/healthcoach-export.service.j2
[Unit]
Description=Health Coach CSV Export to Google Drive
After=network.target

[Service]
Type=oneshot
User={{ ansible_user }}
WorkingDirectory=/opt/healthcoach
EnvironmentFile=/opt/healthcoach/.env
ExecStart=/opt/healthcoach/venv/bin/python3 /opt/healthcoach/app/sync/daily_export.py

[Install]
WantedBy=multi-user.target
