# infra/ansible/roles/healthcoach/templates/healthcoach-decision.service.j2
[Unit]
Description=Health Coach Decision Engine (Daily Plan)
After=docker.service network.target

[Service]
Type=oneshot
User={{ ansible_user }}
WorkingDirectory=/opt/healthcoach
EnvironmentFile=/opt/healthcoach/.env

# Run data sync before generating the plan
ExecStartPre=/opt/healthcoach/venv/bin/python3 app/sync/main.py

ExecStart=/opt/healthcoach/venv/bin/python3 app/decision_engine/main.py
