# infra/ansible/roles/healthcoach/templates/healthcoach-export.timer.j2
[Unit]
Description=Daily Health Coach CSV Export to Google Drive

[Timer]
OnCalendar=*-*-* 20:00:00
Persistent=true

[Install]
WantedBy=timers.target
