# infra/ansible/roles/healthcoach/templates/healthcoach-bot.service.j2
[Unit]
Description=Health Coach Telegram Bot
After=network.target

[Service]
User=cameron
WorkingDirectory=/opt/healthcoach
# This line tells Python that /opt/healthcoach is the root of your code
Environment="PYTHONPATH=/opt/healthcoach"
EnvironmentFile=/opt/healthcoach/.env
ExecStart=/opt/healthcoach/venv/bin/python3 /opt/healthcoach/app/bot/main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
